# 04_knowledge_base/concepts/attacks/ATK-PERSISTENCE-001.yaml
# GR Atom: Persistence
# Version: 1.0
# Created: 2025-01-29

identity:
  id: "ATK-PERSISTENCE-001"
  name: "Persistence"
  normalization:
    normalized_name: "persistence backdoor maintain access"
    normalization_version: "1.0"
  aliases:
    - "Persistence"
    - "지속성"
    - "Maintain Access"
    - "Backdoor"

classification:
  domain: security
  type: concept
  is_infrastructure: false
  abstraction_level: 3
  gr_coordinates:
    layer: "L7"
    zone: "Z3"
    atom_tags: ["CONCEPT", "ATTACK", "PERSISTENCE"]

definition:
  what: |
    시스템 재부팅이나 자격 증명 변경 후에도 접근을
    유지하기 위한 공격 기법 분류이다.
  why: |
    장기적 접근을 유지한다.
    재침투 필요성을 제거한다.
    은밀한 활동을 지속한다.
  how: |
    1. 지속성 메커니즘 선택
    2. 백도어 설치
    3. 은닉
    4. 검증

relations:
  structural:
    is_a:
      - "GR-SEC-ATK-001"
    has_children:
      - "ATK-MITRE-T1547-001"

properties:
  technical:
    techniques:
      windows:
        - "Registry Run Keys"
        - "Scheduled Tasks"
        - "Services"
        - "DLL Hijacking"
        - "WMI Subscription"
      linux:
        - "Cron Jobs"
        - "SSH Keys"
        - "Systemd Services"
        - "Bashrc/Profile"
      web:
        - "Web Shell"
        - "Backdoor Account"
    mitre_tactic: "TA0003"
    detection:
      - "Autoruns 분석"
      - "레지스트리 모니터링"
      - "파일 무결성 검사"

metadata:
  trust:
    source: official
    references:
      - "MITRE ATT&CK TA0003"
    confidence: 1.0
    verified:
      status: verified
      date: 2025-01-29
      by: "GR System"
  temporal:
    created: 2025-01-29
    modified: 2025-01-29
    revision: 1
  security:
    sensitivity: public
  ai:
    embedding_text: |
      Persistence는 시스템 접근을 장기적으로 유지하기 위한
      공격 기법 분류이다.
    search_keywords:
      - "Persistence"
      - "지속성"
      - "Backdoor"
      - "Maintain Access"
    related_queries:
      - "Persistence이란 무엇인가?"
      - "Persistence의 핵심 원리는?"
      - "Persistence을 어떻게 적용하는가?"
